# roles/ansible_user/templates/ansible-pull.timer.j2
[Unit]
Description=Ansible-pull user timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=5m
RandomizedDelaySec=30
Service=ansible-pull.service

[Install]
WantedBy=multi-user.target
